# app/server/itc-os-presence.service
[Unit]
Description=ITC-OS 2026 Presence Server
After=network.target

[Service]
Type=simple
User=rathpisey
WorkingDirectory=/home/rathpisey/ITC-OS-2026/app/server
ExecStartPre=/usr/bin/python3 -m pip install --quiet --break-system-packages -r requirements.txt
ExecStart=/usr/bin/python3 app.py --host 0.0.0.0 --port 5000
Restart=always
RestartSec=3
Environment=TZ=Asia/Phnom_Penh

[Install]
WantedBy=multi-user.target
